# api/deploy/ideam-api.service
[Unit]
Description=IDEAM API (FastAPI sobre PostgreSQL/TimescaleDB)
After=docker.service network-online.target
Wants=network-online.target

[Service]
# Usuario sin privilegios (hardening 2026-06-06): dueño de /var/lib/ideam-api.
User=ideam-api
NoNewPrivileges=yes
PrivateTmp=yes
EnvironmentFile=/etc/ideam/ideam.env
WorkingDirectory=/opt/ideam/api
ExecStart=/opt/ideam/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
TimeoutStopSec=30
# Techo de memoria (auditoria 2026-07-01): el box A1 (24 GB) se comparte con
# Postgres/TimescaleDB (el proceso mas grande y el que mato el OOM-killer en el
# incidente del 2026-06-13) y con los jobs pesados (slice ideam-heavy, 4G).
# MemoryHigh frena a la API (reclaim/throttle del kernel) antes de comprometer
# la RAM del resto; MemoryMax es el limite duro: si un export o pico de
# peticiones concurrentes se desboca, el OOM cae DENTRO de este servicio (que
# se reinicia solo via Restart=always) y no sobre Postgres. 2 workers uvicorn
# operan normalmente muy por debajo de 3G.
MemoryHigh=3G
MemoryMax=4G

[Install]
WantedBy=multi-user.target
