# api/deploy/ideam-api.service
[Unit]
Description=IDEAM API (FastAPI sobre PostgreSQL/TimescaleDB)
After=docker.service network-online.target
Wants=network-online.target

[Service]
User=root
EnvironmentFile=/etc/ideam/ideam.env
WorkingDirectory=/opt/ideam/api
ExecStart=/opt/ideam/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
