# api/deploy/ideam-export-cleanup.timer
[Unit]
Description=Limpieza periodica de exports IDEAM

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
