# deploy/ideam-idf-refresh.service
[Unit]
Description=Refresco mensual de curvas IDF (maximos anuales moviles)
After=docker.service

[Service]
Type=oneshot
Nice=18
IOSchedulingClass=idle
# Serializado con el resto de trabajos pesados via flock compartido (un solo
# recalculo a la vez) y agrupado en el slice que cede recursos a la API.
Slice=ideam-heavy.slice
ExecStart=/usr/bin/flock -w 7200 /run/ideam-heavy.lock /usr/local/bin/ideam-idf-refresh.sh
# Cota explicita: hasta 2h esperando el lock + el refresco (itera estaciones).
TimeoutStartSec=3h
