# deploy/ideam-notify-fail@.service
# ============================================================
# ideam-notify-fail@.service   (unidad PLANTILLA, opcional)
# ------------------------------------------------------------
# Notifica un fallo a healthchecks.io. Es generica: el %i es el
# slug del check a "pingear" como fallido. Se invoca desde el
# OnFailure= de otra unidad, p.ej.:
#     OnFailure=ideam-notify-fail@delta.service
# y eso hace curl a $HEALTHCHECK_URL/fail al fallar ideam-delta.
#
# Requiere HEALTHCHECK_URL en /etc/ideam/ideam.env, p.ej.:
#     HEALTHCHECK_URL=https://hc-ping.com/<uuid-del-check>
# Si la variable no existe, el curl no se ejecuta (no rompe).
# El %i del template es informativo (queda en el journal); el
# ping de fallo va a $HEALTHCHECK_URL/fail.
#
# Crea el check GRATIS en https://healthchecks.io (plan free:
# 20 checks). Copia su ping URL a HEALTHCHECK_URL.
#
# INSTALACION:
#   sudo install -m 644 deploy/ideam-notify-fail@.service /etc/systemd/system/
#   sudo systemctl daemon-reload
# (no se habilita; se referencia desde OnFailure= de otras units)
# ============================================================
[Unit]
Description=Notifica fallo de %i a healthchecks.io

[Service]
Type=oneshot
EnvironmentFile=-/etc/ideam/ideam.env
# Solo hace ping si HEALTHCHECK_URL esta definida. El -m 10 evita
# que un endpoint lento bloquee. /fail marca el check como caido.
ExecStart=/bin/sh -c '[ -n "$HEALTHCHECK_URL" ] && curl -fsS -m 10 --retry 3 "$HEALTHCHECK_URL/fail" || true'
