# deploy/ideam-backfill.service
# ============================================================
# ideam-backfill.service
# ------------------------------------------------------------
# Carga el historico completo Socrata -> Postgres. Es reanudable
# (tabla ingest_state: los pares dataset/anio en 'done' se
# saltan) por lo que reiniciarlo es seguro.
#
# --mode year: baja anio-a-anio en orden descendente y comprime
# cada anio al completarse, de modo que el pico de disco se
# mantiene acotado (nunca se insertan filas en chunks ya
# comprimidos). Es el modo usado en produccion.
#
# Restart=on-failure + RestartSec=60: si el proceso muere con
# error, systemd lo relanza. NO cubre el caso "vivo pero
# colgado": de eso se encarga ideam-watchdog.timer.
# TimeoutStartSec=infinity: el backfill puede tardar dias.
#
# INSTALACION:
#   sudo install -m 644 deploy/ideam-backfill.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl start ideam-backfill      # arranque manual (no enable)
# Seguimiento:  journalctl -u ideam-backfill -f
# ============================================================
[Unit]
Description=IDEAM backfill historico (Socrata -> Postgres, reanudable)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=/etc/ideam/ideam.env
# Reanudable: los (dataset, anio) en estado 'done' se saltan automaticamente.
ExecStart=/opt/ideam/venv/bin/python -m ideam_socrata.db.backfill --dataset all --mode year --compress
WorkingDirectory=/opt/ideam/app
Restart=on-failure
RestartSec=60
TimeoutStartSec=infinity
