# deploy/ideam-fiabilidad.service
[Unit]
Description=Recalculo mensual del semaforo de fiabilidad por estacion (Lote 2.1)
After=docker.service ideam-idf-refresh.service
Wants=ideam-idf-refresh.service

[Service]
Type=oneshot
Nice=18
IOSchedulingClass=idle
# Serializado con el resto de trabajos pesados via flock compartido y agrupado
# en el slice que cede recursos a la API.
Slice=ideam-heavy.slice
User=root
EnvironmentFile=/etc/ideam/ideam.env
WorkingDirectory=/opt/ideam/api
ExecStart=/usr/bin/flock -w 7200 /run/ideam-heavy.lock /opt/ideam/venv/bin/python -m app.services.fiabilidad_batch
# 4h = hasta 2h esperando el lock + el recalculo (que de por si es <=30min).
TimeoutStartSec=4h
