# deploy/ideam-idf-refresh.timer
[Unit]
Description=Refresco IDF: dia 2 de cada mes 06:00 UTC (tras delta 04:00 + backup 05:00)

[Timer]
# 06:00 UTC para correr DESPUES del delta diario (04:00) y del backup (05:00).
# Antes estaba a las 03:30 — ANTES del delta — y recalculaba sobre datos viejos.
OnCalendar=*-*-02 06:00:00 UTC
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
