# deploy/ideam-delta.service
[Unit]
Description=IDEAM delta diario (Socrata -> Postgres)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=/etc/ideam/ideam.env
# Slice + flock compartido: cede CPU/IO a Postgres+API y serializa con los
# demas trabajos pesados (un solo recalculo/ingesta pesada a la vez).
Slice=ideam-heavy.slice
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
# flock -w 7200: si otro trabajo pesado tiene el lock, espera hasta 2h y luego
# corre; el lock se libera solo cuando el otro proceso termina (incl. si muere).
ExecStart=/usr/bin/flock -w 7200 /run/ideam-heavy.lock /opt/ideam/venv/bin/python -m ideam_socrata.db.delta
WorkingDirectory=/opt/ideam/app
# 6h = hasta 2h esperando el lock + el propio delta (que normalmente es corto).
TimeoutStartSec=6h
