# deploy/ideam-delta.timer
[Unit]
Description=Ejecuta el delta IDEAM todos los dias a las 04:00 (hora del servidor)

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
